# systemd/pcscd-restart.timer
[Unit]
Description=Timer to restart pcscd service after boot and login

[Timer]
OnBootSec=1min

[Install]
WantedBy=multi-user.target
